# dep/proasense-reverse-proxy/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "service apache2 start"
